# scripts/go-asr-bot.service
[Unit]
Description=Go ASR Bot - Telegram bot for ASR transcription
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=go-asr-bot
Group=go-asr-bot
WorkingDirectory=/opt/go-asr-bot
ExecStart=/usr/local/bin/go-asr-bot
Restart=always
RestartSec=5
EnvironmentFile=/opt/go-asr-bot/.env

[Install]
WantedBy=multi-user.target
